# files/system/usr/lib/systemd/user/flatpak-bazaar-cleanup.service
[Unit]
Description=Flatpak/Bazaar cleanup on session start
After=default.target

[Service]
Type=oneshot
ExecStart=/usr/bin/bash /usr/libexec/bog/cleanup-flatpak-state.sh

[Install]
WantedBy=default.target