# files/system/etc/profile.d/aliases.sh
# odd shorthand commands that use the system-installed cli tools like ffmpeg and yt-dlp

# getmp4() { yt-dlp }
# getmp3() { yt-dlp }
testfunc() { echo "haiiiii" }