# files/system/etc/profile.d/appimage-compat.sh
#!/usr/bin/env bash
export APPIMAGE_EXTRACT_AND_RUN=1